# scripts/rpi-showcase.service
[Unit]
Description=RPI Demo 3 Showcase PWA
After=network.target graphical-session.target
Wants=graphical-session.target

[Service]
Type=simple
User=pi
Group=pi
Environment=DISPLAY=:0
Environment=XAUTHORITY=/home/pi/.Xauthority
WorkingDirectory=/home/pi/rpi-demo-3
ExecStartPre=/bin/bash -c 'sleep 15'
ExecStart=/home/pi/rpi-demo-3/scripts/kiosk.sh
Restart=always
RestartSec=15
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=graphical.target
