# scripts/rpi-showcase.service
[Unit]
Description=RPI Demo 3 Showcase PWA
After=network.target multi-user.target
Wants=network.target

[Service]
Type=simple
User=pi
Group=pi
Environment=DISPLAY=:0
Environment=XAUTHORITY=/home/pi/.Xauthority
Environment=HOME=/home/pi
WorkingDirectory=/home/pi/rpi-demo-3
ExecStart=/home/pi/rpi-demo-3/scripts/kiosk.sh
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
